# templates/systemd_rc-local.service.epp
[Unit]
 Description=rc.local Compatibility
 ConditionPathExists=<%= $rclocal::config_file %>

[Service]
 Type=forking
 ExecStart=<%= $rclocal::config_file %> start
 TimeoutSec=0
 StandardOutput=tty
 RemainAfterExit=yes
 SysVStartPriority=99

[Install]
 WantedBy=multi-user.target
